# tools/docker/cms/supervisor-app.ini
[program:${docker.project.name}]
command = ${docker.project.path}/env/bin/uwsgi --ini ${docker.project.path}/env/uwsgi.ini
stderr_logfile = /var/log/supervisor/${docker.project.name}-stderr.log
stdout_logfile = /var/log/supervisor/${docker.project.name}-stdout.log

[program:nginx]
command = /sbin/nginx -g 'daemon off;'
stderr_logfile = /var/log/supervisor/nginx-stderr.log
stdout_logfile = /var/log/supervisor/nginx-stdout.log